# superv.conf
[program:rails-blog]
command=rails server
directory=/var/www/railes_guides
autostart=true
autorestart=true